# config/supervisord.conf

[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor.log
logfile_maxbytes=1MB
logfile_backups=2
loglevel=debug
nodaemon=false
directory=%(here)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:clamav_init]
command=/opt/start_clamav.sh
autostart=true
autorestart=false
priority=10
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:clamd]
command=/usr/sbin/clamd -c /etc/clamav/clamd.conf
autostart=false
autorestart=true
priority=30
user=clamav
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log


[program:freshclam]
command=/usr/bin/freshclam -d --quiet -c %(ENV_FRESHCLAM_UPDATE)s --config-file=/etc/clamav/freshclam.conf
autostart=false
autorestart=true
priority=20
startretries=1
user=clamav
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
